# config/cron/rootful.cron
### BEGIN CONTAINERS CRON JOBS ###
10 5 1 * * /bin/podman image prune -af
### END CONTAINERS CRON JOBS ###
