# config/cron/rootless.cron
### BEGIN CONTAINERS CRON JOBS ###
5 3 * * * /bin/bash {REPO_DIR}/bin/backup.sh
5 5 1 * * /usr/bin/podman image prune -af
### END CONTAINERS CRON JOBS ###
