# config/cron/rootful.cron
### BEGIN CONTAINERS CRON JOBS ###
0 3 * * * /bin/bash {REPO_DIR}/bin/backup.sh
0 5 1 * * /bin/podman image prune -af
### END CONTAINERS CRON JOBS ###
