# rootful/restic/restic-root.timer
[Unit]
Description=Run restic-root daily

[Timer]
OnCalendar=*-*-* 03:20:00
Persistent=true

[Install]
WantedBy=timers.target
