# rootful/restic/restic-home-share.timer
[Unit]
Description=Run restic-home-share daily

[Timer]
OnCalendar=*-*-* 03:40:00
Persistent=true

[Install]
WantedBy=timers.target
